# etc/supervisord.d/hbase-regionserver.ini
[program:hbase-regionserver]
command = /opt/qnib/hbase/regionserver/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
