# etc/supervisord.d/hbase-master.ini
[program:hbase-master]
command = /opt/qnib/hbase/master/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
